# ops/systemd/download_metrics.timer
[Unit]
Description=Run metric fetcher every day

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target